# packages/etc/inittab
::sysinit:/os/setup.sh
tty1::respawn:/os/getty 38400 tty1
