# packages/etc/inittab
::sysinit:/os/mount -a
::respawn:/os/getty -l /os/login 115200 ttyS0
